# tesis/services/minoir.service
[Unit]
After=vpn.service

[Service]
ExecStart=/etc/systemd/system/minoir.sh

[Install]
WantedBy=default.target